# config/supervisor_taurus.conf
[supervisord]
nodaemon=true

[program:wsbitfinex]
command=/usr/bin/php artisan taurus:websocket_bitfinex
directory=/root
startretries=3
stopwaitsecs=10
autostart=true